# .config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.3.22)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/screaminginc/.nvm/versions/node/v22.22.1/bin/node /mnt/c/Users/SCREAMING.INC/AppData/Roaming/npm/node_modules/openclaw/dist/index.js gateway --port 18789
Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group
Environment=HOME=/home/screaminginc
Environment=TMPDIR=/tmp
Environment=NODE_EXTRA_CA_CERTS=/etc/ssl/certs/ca-certificates.crt
Environment=PATH=/home/screaminginc/.nvm/versions/node/v22.22.1/bin:/home/screaminginc/.nvm/current/bin:/home/screaminginc/.local/bin:/home/screaminginc/.npm-global/bin:/home/screaminginc/bin:/home/screaminginc/.volta/bin:/home/screaminginc/.asdf/shims:/home/screaminginc/.bun/bin:/home/screaminginc/.fnm/current/bin:/home/screaminginc/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Gateway"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.3.22

[Install]
WantedBy=default.target
